# deploy/com.wesleyflorence.recall.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.wesleyflorence.recall</string>

    <key>ProgramArguments</key>
    <array>
      <string>/Users/wesleyflorence/.asdf/installs/bun/1.3.6/bin/bun</string>
      <string>run</string>
      <string>start</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/wesleyflorence/code/recall</string>

    <key>EnvironmentVariables</key>
    <dict>
      <key>PATH</key>
      <string>/Users/wesleyflorence/.asdf/installs/bun/1.3.6/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
      <key>PORT</key>
      <string>18104</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <dict>
      <key>SuccessfulExit</key>
      <false/>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/wesleyflorence/Library/Logs/recall/recall.stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/wesleyflorence/Library/Logs/recall/recall.stderr.log</string>
  </dict>
</plist>
