# deploy/com.wesleyflorence.recall.backup.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.wesleyflorence.recall.backup</string>

    <key>ProgramArguments</key>
    <array>
      <string>/bin/bash</string>
      <string>/Users/wesleyflorence/code/recall/deploy/backup-recall-postgres.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/wesleyflorence/code/recall</string>

    <key>StartCalendarInterval</key>
    <dict>
      <key>Hour</key>
      <integer>3</integer>
      <key>Minute</key>
      <integer>0</integer>
    </dict>

    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>/Users/wesleyflorence/Library/Logs/recall/recall-backup.stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/wesleyflorence/Library/Logs/recall/recall-backup.stderr.log</string>
  </dict>
</plist>
